# roles/after_started_systemd/templates/ssh-agent-npiperelay@.service.j2
[Unit]
Description=ssh-anget relay service

[Service]
ExecStart=%h/.local/bin/npiperelay.exe -ei -s -v //./pipe/openssh-ssh-agent
StandardInput=socket
StandardOutput=socket
StandardError=journal
